# dotfiles/systemd/user/delete-old-ss.timer
[Unit]
Description=run delete old screenshots monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
