# dotfiles/systemd/user/dunst.service
[Unit]
Description=Dunst notification daemon
PartOf=xsession.target

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=xsession.target
